# 00-usb-picture-frame.rules
ACTION=="add", SUBSYSTEM=="usb_device", ATTRS{idVendor}=="1b8e", ATTRS{idProduct}=="0021", OWNER="pi", RUN+="/home/pi/photo-frame/copy-pictures.sh"
