# 00-usb-picture-frame.rules
ACTION=="add", ATTRS{idVendor}=="1b8e", ATTRS{idProduct}=="0021", RUN+="/home/pi/picture-frame/copy-pictures.sh %k"
